# scripts/root-resume.service
[Unit]
Description=Local system resume actions
After=suspend.target

[Service]
Type=simple
ExecStart=/bin/chvt 1

[Install]
WantedBy=suspend.target
